# dns-watchdog.service
[Unit]
Description=DNS Watchdog Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dns-watchdog.sh
